# linux-3.18.29-dev/root_fs/chaos_calmer/etc/init.d/setup
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=15
boot() {
    [ -f /proc/irq/234/smp_affinity ] && echo "8" > /proc/irq/234/smp_affinity
}
