# chaos_calmer_15_05_1/target/linux/siflower/sf16a18-fullmask/base-files-SF16A18-REP-NOPA/etc/init.d/system
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=10
USE_PROCD=1

validate_system_section()
{
	uci_validate_section system system "${1}" \
		'hostname:string:OpenWrt' \
		'conloglevel:uinteger' \
		'buffersize:uinteger' \
		'timezone:string:UTC' \
		'zonename:string'
}


system_config() {
	local cfg="$1"

	local hostname conloglevel buffersize timezone zonename

	validate_system_section "${1}" || {
		echo "validation failed"
		return 1
	}
#add mac for suffix
	hostnameset=`uci -q get system.@system[0].hostnameset`
	[ -z "$hostnameset" ] && {
		mac=`uci get network.wan.macaddr`
		[ -z "$mac"  ] && mac=`uci get network.lan.macaddr`
		hostname=$hostname${mac:12:2}${mac:15:2}
		eval "uci -q set system.@system[0].hostname=${hostname}"
		eval "uci -q set system.@system[0].hostnameset=1"
		uci commit system
	}

	echo "$hostname" > /proc/sys/kernel/hostname
	[ -z "$conloglevel" -a -z "$buffersize" ] || dmesg ${conloglevel:+-n $conloglevel} ${buffersize:+-s $buffersize}
	echo "$timezone" > /tmp/TZ
	[ -n "$zonename" ] && [ -f "/usr/share/zoneinfo/$zonename" ] && ln -s "/usr/share/zoneinfo/$zonename" /tmp/localtime

	# apply timezone to kernel
	date -k
}

reload_service() {
	config_load system
	config_foreach system_config system
}

service_triggers()
{
	procd_add_reload_trigger "system"
	procd_add_validation validate_system_section
}

start_service() {
	reload_service
}
