# chaos_calmer_15_05_1/target/linux/siflower/sf16a18-fullmask/base-files-SF16A18-86V-C2/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
. /lib/functions/system.sh
start_time=`date +%s`
uci set system.ntp.start=$start_time
uci commit
exit 0
