# linux-3.18.29-dev/root_fs/root-siflower-wifi/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/etc/banner.failsafe
/etc/config/samba
/etc/resolv.conf
/etc/init.d/subservice
/etc/init.d/done
/bin/config_generate
/var
/etc/opkg/keys/de98a2dd1d0f8a07
/etc/config/siserver
/etc/init.d/system
/etc/preinit
/lib/upgrade/platform.sh
/etc/shells
/lib/preinit/10_sysinfo
/etc/board.d/01_network
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/lib/preinit/10_indicate_preinit
/lib/functions/uci-defaults-new.sh
/sbin/insert_sd.sh
/etc/config/system
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/etc/opkg/keys/53bad1233d4c98c5
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/config/basic_setting
/etc/shadow
/etc/config/qos_cfg
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/10_indicate_failsafe
/init
/etc/opkg/keys/cd7de88535e37b08
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/usr/bin/otaupgrade
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/config/sibmob
/etc/config/fstab
/etc/rc.button/failsafe
/etc/services
/lib/functions/system.sh
/bin/board_detect
/etc/sysupgrade_except.conf
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/etc/rc.button/power
/etc/mdev.conf
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/config/notify
/lib/preinit/40_run_failsafe_hook
/usr/bin/getSN
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/sbin/firstboot
/usr/bin/init_devlist
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/functions.sh
/etc/config/devlist
/sbin/sd_card_automount
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/etc/profile
/etc/uci-defaults/02_network
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/usr/bin/speed
/lib/functions/service.sh
/etc/init.d/umount
/etc/init.d/syncservice
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/lib/preinit/02_default_set_state
/etc/mtab
/etc/init.d/boot
/etc/config/device_bomb
/lib/functions/preinit.sh
/etc/config/siwifi
/bin/ipcalc.sh
/etc/diag.sh
