# linux-3.18.29-dev/root_fs/root-siflower-wifi/usr/lib/opkg/info/samba36-server.list
/usr/sbin/nmbd
/etc/samba/smb.conf.template
/etc/samba/valid.dat
/etc/init.d/samba
/etc/samba/upcase.dat
/usr/sbin/smbpasswd
/usr/sbin/smbd
/usr/sbin/samba_multicall
/etc/samba/lowcase.dat
