# chaos_calmer_15_05_1/target/linux/siflower/sf16a18-fullmask/base-files-SF16A18-86V-C2/etc/init.d/led_86v
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=97


start() {
	echo 1 > /tmp/wifi_onoff_status
	/bin/led-button -l 25
	./etc/anomaly_server &
}

