# linux-3.18.29-dev/root_fs/ltp_test/etc/inittab
::sysinit:/etc/init.d/rcS
ttyS0::askfirst:/bin/ash --login
::ctrlaltdel:/bin/umount -a -r
